# ups-monitor-v2.service
[Unit]
Description=X120x UPS Monitoring (v2 canary, no shutdown)
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=/home/pi/x120x
ExecStart=/usr/bin/python3 -u /home/pi/x120x/merged_v2.py --dry-run --no-pld
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
